# app/homeassistant/openwrt/rc.local
#!/bin/sh

# hostapd to mqtt hook setup
mkdir -p /var/lib/home-assistant
config_file="/etc/hostapd_hook_mqtt.conf"
if [ ! -e "$config_file"; then
  echo "error: could not start hostapd_cli for mqtt hook: config file $config_file not found"
else
  echo "start hostapd_cli as daemon to call /bin/hostapd_mqtt_hook.sh on events"
  hostapd_cli -i global -B -r -a /bin/hostapd_hook_mqtt.sh
fi
