# appliance/systemd/mail-to-sentry.path
[Unit]
Description=watch for /var/mail/root/new changes

[Path]
PathChanged=/var/mail/root/new

[Install]
WantedBy=prepare-appliance.service
