# app/gitops/gitops-update.timer
[Unit]
Description=Gitops Update Timer

[Timer]
OnCalendar={{ settings.timer.oncalendar }}

[Install]
WantedBy=timers.target
