# containers/template/podman-system-prune.service
[Unit]
Description=Force a podman system prune

[Service]
Type=oneshot
ExecStart=podman system prune -f
