# app/backup/backup.timer
[Unit]
Description=Backup Timer

[Timer]
OnCalendar={{ settings.timer }}

[Install]
WantedBy=timers.target
